# rrc-inventory-mdns.service
[Unit]
Description=Publish rrc-inventory.local via mDNS (Avahi)
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/srinath/RRC-Inventory/scripts/mdns_publish.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
